# deploy/agents/crontab
# Бомбора Agent Pipeline (время по Батуми, Asia/Tbilisi UTC+4)
# Все 4 агента последовательно на одной ветке, один PR в конце
#
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Утренний прогон — 11:00
0 11 * * * su -s /bin/bash agent -c '/scripts/run-pipeline.sh >> /logs/agents.log 2>&1'

# Вечерний прогон — 19:00
0 19 * * * su -s /bin/bash agent -c '/scripts/run-pipeline.sh >> /logs/agents.log 2>&1'
