# deploy/agents/crontab
# Бомбора Agent Schedule (время по Батуми, Asia/Tbilisi UTC+4)
# Logs: /logs/{agent}.log
#
# Product Agent — каждый день в 10:00
0 10 * * * /scripts/run-agent.sh product >> /logs/agents.log 2>&1

# Designer Agent — каждый день в 12:00
0 12 * * * /scripts/run-agent.sh designer >> /logs/agents.log 2>&1

# Developer Agent — каждый день в 16:00
0 16 * * * /scripts/run-agent.sh developer >> /logs/agents.log 2>&1

# Tech Lead Agent — каждый день в 20:00
0 20 * * * /scripts/run-agent.sh tech-lead >> /logs/agents.log 2>&1
